# docker/supervisord/wecamp_bot.conf
[program:wecamp_bot]
command=/usr/bin/php bot.php
directory=/var/bot
stderr_logfile=/var/bot/var/errorlog.txt
startsecs=3
autorestart=true
